# ops/systemd/turbovas-quality-gate.timer.in
# SPDX-FileCopyrightText: 2026 Robert Pelfrey <Robert@Pelfrey.de>
# SPDX-License-Identifier: GPL-3.0-or-later

[Unit]
Description=Run the TurboVAS local quality gate daily
Documentation=file:@REPO_ROOT@/README.md

[Timer]
OnCalendar=*-*-* 03:30:00
Persistent=true
RandomizedDelaySec=20m
Unit=turbovas-quality-gate.service

[Install]
WantedBy=timers.target
