# ops/systemd/yafvs-quality-gate.service.in
# SPDX-FileCopyrightText: 2026 Robert Pelfrey <robert@pelfrey.de>
# SPDX-License-Identifier: GPL-3.0-or-later

[Unit]
Description=YAFVS local quality gate
Documentation=file:@REPO_ROOT@/README.md

[Service]
Type=oneshot
WorkingDirectory=@REPO_ROOT@
ExecStart=@PYTHON@ tools/yafvsctl quality-gate --json
StandardOutput=append:@QUALITY_GATE_LOG_DIR@/quality-gate.log
StandardError=append:@QUALITY_GATE_LOG_DIR@/quality-gate.log
